# crontab
* * * * *  /usr/bin/php /var/www/laravel/storage/updateDB.php

